# lib/bulletin-notify/bulletinnotify.conf
[program:bulletinnotify]
directory=/path/to/bulletin-notify
command=/path/to/venv/bin/python bulletinnotify.py
autostart=true
autorestart=true
stdout_logfile=/var/log/bulletin/bulletinnotify.out.log
stderr_logfile=/var/log/bulletin/bulletinnotify.error.log
